# pi-control.service
[Unit]
Description=Pi Control Panel API
Documentation=https://github.com/BGirginn/rasp_pi_webUI
After=network.target

[Service]
Type=simple
User=bgirgin
Group=bgirgin
WorkingDirectory=/opt/pi-control/panel/api
Environment=PYTHONPATH=/opt/pi-control/panel/api
Environment=DATABASE_PATH=/var/lib/pi-control/control.db
Environment=TELEMETRY_DB_PATH=/var/lib/pi-control/telemetry.db
Environment=AGENT_SOCKET=/run/pi-control/agent.sock
Environment=JWT_SECRET_FILE=/etc/pi-control/jwt_secret
Environment=API_DEBUG=false
ExecStart=/opt/pi-control/venv/bin/uvicorn main:app --host 127.0.0.1 --port 8080
Restart=always
RestartSec=5
TimeoutStopSec=5
KillMode=mixed

# Security hardening (relaxed for terminal PTY access)
ProtectSystem=full
ProtectHome=false
ReadWritePaths=/var/lib/pi-control /opt/pi-control /home/bgirgin
PrivateTmp=true

[Install]
WantedBy=multi-user.target
