# agent/pi-agent.service
[Unit]
Description=Pi Control Panel Agent
Documentation=https://github.com/bgirginn/rasp_pi_webUI
After=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/pi-control/agent

# Main process
ExecStart=/opt/pi-control/venv/bin/python3 /opt/pi-control/agent/pi-agent.py
ExecReload=/bin/kill -HUP $MAINPID

# Restart policy
Restart=on-failure
RestartSec=10s
StartLimitInterval=5min
StartLimitBurst=5

# Watchdog
WatchdogSec=60

# Environment
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=-/opt/pi-control/.env

# Minimal Security/Permissions (Permissive for troubleshooting)
# NoNewPrivileges=false
# PrivateTmp=false
# ProtectSystem=false

# Socket directory
RuntimeDirectory=pi-agent
RuntimeDirectoryMode=0755

# Capabilities
CapabilityBoundingSet=CAP_NET_ADMIN CAP_SYS_PTRACE
AmbientCapabilities=CAP_NET_ADMIN CAP_SYS_PTRACE

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pi-agent

[Install]
WantedBy=multi-user.target
